# docker/src/backupDB.cron
01 * * * * /root/backupDB.sh > /var/log/cron 2>&1
