# docker/src/backupDB.cron
01 23 * * * root /root/backupDB.sh > /dev/null 2>&1
